# baseboard-bleedingedge/chromeos-base/chromeos-bsp-baseboard-bleedingedge/files/init/load-rtd1619-modules.conf
#  Copyright 2024 The FydeOS Authors. All rights reserved.
#  Use of this source code is governed by a BSD-style license that can be
#  found in the LICENSE file.
#  Author: Yang Tsao<yang@fydeos.io>

description "Job to load modules"
author      "yang@fydeos.io"

oom score never
task

start on starting udev-trigger-early

script
/usr/sbin/loadmodules.sh
end script
