# baseboard-bleedingedge/sys-apps/firmware-log-tools/files/rkfwlog-init.conf
# Copyright 2024 The fydeOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Setting firmware log level when boot-services is complete"
author         "yang@fydeos.io"

env RKFWLOG_LEVEL=0

task
oom score never

start on started boot-services

script
  modprobe -q rtk_fwdbg
  modprobe -q rtk_avcpulog
  if [ -f /usr/local/bin/fwdbg ]; then
    /usr/local/bin/fwdbg ${RKFWLOG_LEVEL}
  else
    /usr/bin/fwdbg ${RKFWLOG_LEVEL}
  fi
end script
